# cookbooks/apitrary_pytools/templates/default/supervisor_deployr.conf.erb
[program:deployr]
command=/usr/local/bin/celery worker -A pydeployr.deployr_tasks -l info
user=root
directory=/home/deployr
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/deployr.log
loglevel=info