# cookbooks/apitrary_pytools/templates/default/supervisor_deployr.conf.erb
[program:deployr_<%= @deployr_mode%>]
command=/usr/bin/python /usr/local/bin/deployr.py --mode=<%= @deployr_mode%>
directory=/home/deployr
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/deployr_<%= @deployr_mode%>.log
loglevel=info