# cookbooks/apitrary_pytools/templates/default/supervisor_trackr.conf.erb
[program:trackr]
command=/usr/bin/python /usr/bin/trackr.py
directory=/home/trackr
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/trackr.log
loglevel=info