# cookbooks/apitrary_pytools/templates/default/supervisor_balancr.conf.erb
[program:balancr]
command=/usr/local/bin/celery worker -A pytools.pybalancr.balancr_tasks -l info
directory=/home/balancr
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/balancr.log