# cookbooks/apitrary_pytools/templates/default/supervisor_balancr.conf.erb
[program:balancr]
command=/usr/bin/python /usr/local/bin/deployr.py --mode=balance
directory=/home/balancr
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/balancr.log