# cookbooks/apitrary_myip/templates/default/supervisor_myip.conf.erb
[program:myip]
command=/usr/bin/python /usr/local/bin/myip_runner.py
user=root
directory=/home/myip
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/myip.log
loglevel=info