# cookbooks/apitrary_bugtracking/templates/default/supervisor_youtrack.conf.erb
[program:youtrack]
command=/usr/bin/java -Xmx512m -Djava.awt.headless=true -jar /home/bugs/youtrack/youtrack-4.1.jar 8000
directory=/home/bugs
user=bugs
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/youtrack.log
loglevel=info
startsecs=10
stopsignal=KILL
